# gps-manager.service
[Unit]
Description=reading from CAN looking for telemetry messages

[Install]
WantedBy=multi-user.target

[Service]
ExecStart=/path/to/can_listener
